# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for superlists-staging.superlists.ml"

start on net-device-up
stop on shutdown

respawn

setuid elspeth
chdir /home/elspeth/sites/superlists-staging.superlists.ml/source

exec ../virtualenv/bin/gunicorn \
--bind unix:/tmp/superlists-staging.superlists.ml.socket \
superlists.wsgi:application
